# base/conf/shell/etc/profile.d/00-reset-path.sh
# Reset default PATH settings
PATH="/opt/modular/pkg/packages.modular.com_mojo/bin:/opt/code-server/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
export PATH
